# udev/99-frgb.rules
# frgb — Lian Li USB device permissions
# Install: sudo cp 99-frgb.rules /etc/udev/rules.d/
# Reload:  sudo udevadm control --reload-rules && sudo udevadm trigger

# --- USB device permissions (allow non-root access) ---

# Wireless TX dongle (VID=0416)
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="8040", MODE="0666"
# Wireless RX dongle (VID=0416)
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="8041", MODE="0666"
# Wireless TX dongle V2 (CH340)
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="e304", MODE="0666"
# Wireless RX dongle V2 (CH340)
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="e305", MODE="0666"

# Strimer Plus LED strip
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a200", MODE="0666"

# ENE 6K77 SL/AL fan controllers
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a100", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a101", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a102", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a103", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a104", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a105", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a106", MODE="0666"

# TL Fan controller
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7372", MODE="0666"
# TLLCD
SUBSYSTEM=="usb", ATTR{idVendor}=="04fc", ATTR{idProduct}=="7393", MODE="0666"

# Galahad II Trinity AIO
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7371", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7373", MODE="0666"

# HydroShift LCD AIO
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7398", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7399", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="739a", MODE="0666"

# Galahad II LCD / Vision AIO
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7391", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7395", MODE="0666"

# Universal Screen 8.8" LED Ring
SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="8050", MODE="0666"

# SLV3 wireless LCD fan
SUBSYSTEM=="usb", ATTR{idVendor}=="1cbe", ATTR{idProduct}=="0005", MODE="0666"
# TLV2 wireless LCD fan
SUBSYSTEM=="usb", ATTR{idVendor}=="1cbe", ATTR{idProduct}=="0006", MODE="0666"
# HydroShift II LCD Circle
SUBSYSTEM=="usb", ATTR{idVendor}=="1cbe", ATTR{idProduct}=="a021", MODE="0666"
# HydroShift II LCD Square
SUBSYSTEM=="usb", ATTR{idVendor}=="1cbe", ATTR{idProduct}=="a034", MODE="0666"
# Lancool 207 Digital
SUBSYSTEM=="usb", ATTR{idVendor}=="1cbe", ATTR{idProduct}=="a065", MODE="0666"
# Universal Screen 8.8"
SUBSYSTEM=="usb", ATTR{idVendor}=="1cbe", ATTR{idProduct}=="a088", MODE="0666"

# SLV3H wired hub (CH340, VID=1A86)
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="2107", MODE="0666"
# CH340 USB hubs used by Lian Li controllers
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="8091", MODE="0666"

# Desktop-mode LCD devices (CH340)
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="acd1", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="ace1", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="ad11", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="ad20", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="ad21", MODE="0666"
SUBSYSTEM=="usb", ATTR{idVendor}=="1a86", ATTR{idProduct}=="ad22", MODE="0666"

# --- Force usbhid to bind HID interfaces the kernel may not claim ---

# HydroShift LCD AIO
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7398", RUN+="/bin/sh -c 'echo 0416 7398 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7399", RUN+="/bin/sh -c 'echo 0416 7399 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="739a", RUN+="/bin/sh -c 'echo 0416 739a > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"

# Galahad II LCD / Vision AIO
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7391", RUN+="/bin/sh -c 'echo 0416 7391 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="7395", RUN+="/bin/sh -c 'echo 0416 7395 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"

# Universal Screen 8.8" LED Ring
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0416", ATTR{idProduct}=="8050", RUN+="/bin/sh -c 'echo 0416 8050 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"

# Strimer Plus LED strip
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a200", RUN+="/bin/sh -c 'echo 0cf2 a200 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"

# ENE 6K77 SL/AL fan controllers
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a100", RUN+="/bin/sh -c 'echo 0cf2 a100 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a101", RUN+="/bin/sh -c 'echo 0cf2 a101 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a102", RUN+="/bin/sh -c 'echo 0cf2 a102 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a103", RUN+="/bin/sh -c 'echo 0cf2 a103 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a104", RUN+="/bin/sh -c 'echo 0cf2 a104 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a105", RUN+="/bin/sh -c 'echo 0cf2 a105 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0cf2", ATTR{idProduct}=="a106", RUN+="/bin/sh -c 'echo 0cf2 a106 > /sys/bus/usb/drivers/usbhid/new_id 2>/dev/null || true'"

# --- hidraw permissions (for hidapi access) ---

# Wireless TX dongle
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="8040", MODE="0666"
# Wireless RX dongle
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="8041", MODE="0666"
# Wireless TX dongle V2 (CH340)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="e304", MODE="0666"
# Wireless RX dongle V2 (CH340)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="e305", MODE="0666"

# SLV3H wired hub
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="2107", MODE="0666"

# SL-LCD Wireless
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1cbe", ATTRS{idProduct}=="0005", MODE="0666"
# TLV2 wireless LCD
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1cbe", ATTRS{idProduct}=="0006", MODE="0666"
# HydroShift II Circle
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1cbe", ATTRS{idProduct}=="a021", MODE="0666"
# HydroShift II Square
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1cbe", ATTRS{idProduct}=="a034", MODE="0666"

# TL Fan controller
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7372", MODE="0666"

# Galahad II Trinity AIO
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7371", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7373", MODE="0666"

# HydroShift LCD AIO
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7398", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7399", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="739a", MODE="0666"

# Galahad II LCD / Vision AIO
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7391", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="7395", MODE="0666"

# Universal Screen 8.8" LED Ring
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0416", ATTRS{idProduct}=="8050", MODE="0666"

# Strimer Plus LED strip
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a200", MODE="0666"

# ENE 6K77 SL/AL fan controllers
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a100", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a101", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a102", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a103", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a104", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a105", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="0cf2", ATTRS{idProduct}=="a106", MODE="0666"

# Desktop-mode LCD devices (CH340)
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="acd1", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="ace1", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="ad11", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="ad20", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="ad21", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1a86", ATTRS{idProduct}=="ad22", MODE="0666"

# ---------------------------------------------------------------------------
# Motherboard hwmon PWM access (Super I/O chips — nct6xxx, it87xx, etc.)
#
# Group-gated (not 0666 like the USB rules above) because writing to pwm*
# can stall fans / overheat hardware if misused. User must be in the `frgb`
# group to control motherboard fans.
#
# Install:
#   sudo groupadd -r frgb
#   sudo usermod -aG frgb $USER    # log out/in to take effect
#   sudo cp udev/99-frgb.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules && sudo udevadm trigger -s hwmon
# ---------------------------------------------------------------------------
SUBSYSTEM=="hwmon", KERNEL=="hwmon*", RUN+="/bin/sh -c 'for f in /sys%p/pwm* ; do [ -f $$f ] && chmod 0664 $$f && chgrp frgb $$f 2>/dev/null || true ; done'"
